# sonar/docker/mysql.conf
[supervisord]
nodeamon=true

[program:mysqld]
command=/usr/sbin/mysqld